# helper/autostart/guystation.desktop
[Desktop Entry]
Type=Application
Exec=sudo npm --prefix=~/guystation start
Hidden=false
X-GNOME-Autostart-enabled=true
Name=GuyStation
