# helper/autostart/guystation.desktop
[Desktop Entry]
Type=Application
Exec=bash -c "xhost + local:; sudo npm --prefix=~/guystation start"
Hidden=false
X-GNOME-Autostart-enabled=true
Name=GuyStation
